# launchd/template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.cc-knowledge</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__SCRIPT__</string>
    </array>

    <!--
    Placeholders substituted by setup:
      __PYTHON__   — absolute path to python3 (e.g. /usr/bin/python3)
      __SCRIPT__   — absolute path to installed daily_ingest.py
      __LOG_DIR__  — directory for stdout/stderr logs (e.g. /Users/you/.cc-knowledge/)
      __API_KEY__  — Anthropic API key injected at install time
    -->

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>__API_KEY__</string>
    </dict>

    <!-- Run once per day at 9:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__ingest.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__ingest.log</string>

    <!-- Do not run immediately on launchctl load -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
